# systemd/panel-backup-panel.service
# update: panel-only backup service unit
[Unit]
Description=Remnawave panel backup (optional Telegram upload)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/panel-backup.env
ExecStart=/usr/bin/env BACKUP_INCLUDE_OVERRIDE=all /usr/local/bin/panel-backup.sh
User=root
Group=root
