# systemd/panel-backup.timer
# update: panel backup timer unit
[Unit]
Description=Run panel backup daily

[Timer]
OnCalendar=*-*-* 03:40:00 UTC
Persistent=true
Unit=panel-backup.service

[Install]
WantedBy=timers.target
