# systemd/panel-backup-bedolaga.service
# update: bedolaga-only backup service unit
[Unit]
Description=Bedolaga backup (bot + cabinet, optional Telegram upload)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/panel-backup.env
ExecStart=/usr/bin/env BACKUP_INCLUDE_OVERRIDE=bedolaga /usr/local/bin/panel-backup.sh
User=root
Group=root
